# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:ephesoft]
#command=/ephesoft/init.sh
#user=alfresco
